# docker/configurations/supervisor/players.ini
[supervisord]
nodaemon=true

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php artisan queue:work --queue=players
directory=/var/www/
user=www
autostart=true
autorestart=true
startsecs=0
numprocs=4
