# docker/configurations/supervisor/replays.ini
[supervisord]
nodaemon=true

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --queue=replays
directory=/var/www/
user=www
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
